# supervisord.conf
[supervisord]
nodaemon=true

[program:scheduler]
command=airflow scheduler
autorestart=true

[program:server]
command=airflow webserver -p 8080
autorestart=true
